# default/profile.d/firefox.sh
#!/bin/bash

export MOZ_DISABLE_RDD_SANDBOX=1
export MOZ_X11_EGL=1
export MOZ_WEBRENDER=1
